# ansible/roles/configure-zfs-snapshots/templates/zfs-snapshot@.timer.j2
[Unit]
Description=ZFS auto-snapshot timer ({{ item.frequency }})

[Timer]
OnCalendar={{ item.calendar }}
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
